# services/system-misc/stackdriver-agent/stackdriver-restart-if-down.timer
[Unit]
Description=stackdriver-restart-if-down timer

[Timer]
OnUnitActiveSec=120
OnBootSec=120s

[Install]
WantedBy=timers.target